# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:xltv]
command=sh -x /opt/xltv/command.sh